# scripts/epicscagd-tst
#!/bin/bash
#
#	/etc/rc.d/init.d/epicscagd-tst
#
# Starts the EPICS CAG gateway
#
# chkconfig: - 99 1
# description: EPICS Channel Access Gateway
# processname: epicscagd-tst

# Source function library.
. /etc/init.d/functions
. /etc/init.d/epicscagp

gwname=tst
gwprefix=NET:CAG:TST
gwprocport=30110

export EPICS_CA_ADDR_LIST=$TST_BC
export EPICS_CAS_INTF_ADDR_LIST=`echo $PSCAG4_IFLIST | sed -e "s%${TST_IF4}%%"`
#cagprogram="/reg/g/pcds/epics/extensions/gateway/R2.1.2.0-1.0.0/bin/$EPICS_HOST_ARCH/gateway"

# Add -no_cache to gwargs to disable caching in gateway
# With -no_cache, CA updates are passed along to clients w/o pre-fetching the value.
gwargs="$gwargs -no_cache"

case "$1" in
start)
	echo Running: cagstart $gwname $gwprefix $gwprocport $gwargs
	cagstart $gwname $gwprefix $gwprocport $gwargs
	;;
name)
	echo "Gateway $gwname (prefix $gwprefix) runs on procServ port $gwprocport"
	;;
*)
	echo $"Usage: $0 {start|name}"
	exit 1
esac

exit $?
