# scripts/epicscagd-fee
#!/bin/bash
#
#	/etc/rc.d/init.d/epicscagd-fee
#
# Starts the EPICS CAG daemon
#
# chkconfig: - 99 1
# description: EPICS Channel Access Gateway
# processname: epicscagd-fee

# Source function library.
. /etc/init.d/functions
. /etc/init.d/epicscagd

RETVAL=0

export EPICS_CA_ADDR_LIST=$FEE_BC
export EPICS_CAS_INTF_ADDR_LIST="$MCC_IF $CDS_IF $CXI_IF $DEV_IF $SXR_IF $AMO_IF $XPP_IF $XCS_IF"
export EPICS_CA_MAX_ARRAY_BYTES=36000000

case "$1" in
start)
	cagstart fee NET:CAG:FEE
	;;
stop)
	cagstop fee
	;;
restart)
	cagrestart fee
	;;
reload)
	cagreload NET:CAG:FEE
	;;
*)
	echo $"Usage: $0 {start|stop|restart}"
	exit 1
esac

exit $?

