# scripts/pvagw-las-rfc-wave.sh
#!/bin/bash
#
# Starts the EPICS PVA gateway
#
# chkconfig: - 99 1
# description: EPICS PVA Gateway
# processname: pvagw-las-rfc-wave

# Source function library.
. /etc/init.d/functions
. /etc/init.d/epicscagp

gwname=las-rfc-wave
gwprocport=40098
export gwHostNum=3

case "$1" in
start)
	pvaStart $gwname $gwprocport
	exit $?
	;;
name)
	echo "PVA Gateway $gwname runs on pscag$gwHostNum w/ procServ port $gwprocport"
	;;
env-only)
	;;
*)
	echo $"Usage: $0 {start|name|env-only}"
	return 1
esac

return $?
